# deploy/install/systemd/loki.service
[Unit]
Description=Loki (ongrid log store)
Documentation=https://grafana.com/docs/loki/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ongrid-loki
Group=ongrid
ExecStart=/usr/local/bin/loki \
  -config.file=/etc/ongrid/loki-config.yaml
KillMode=mixed
TimeoutStopSec=30s
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ongrid-loki

NoNewPrivileges=true
ProtectSystem=true
ProtectHome=true
PrivateTmp=true
StateDirectory=ongrid-loki
StateDirectoryMode=0755
ReadWritePaths=/var/lib/ongrid-loki
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
