# deploy/install/systemd/ongrid-frontier.service
[Unit]
Description=ongrid frontier (edge tunnel multiplexer)
Documentation=https://github.com/singchia/frontier
After=network-online.target
Wants=network-online.target
StartLimitInterval=60s
StartLimitBurst=3

[Service]
Type=simple
User=ongrid
Group=ongrid
EnvironmentFile=-/etc/ongrid/ongrid-frontier.env
WorkingDirectory=/var/lib/ongrid
ExecStart=/usr/local/bin/ongrid-frontier --config /etc/ongrid/frontier.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=20s
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ongrid-frontier

NoNewPrivileges=true
ProtectSystem=true
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=false
RemoveIPC=true
StateDirectory=ongrid
LogsDirectory=ongrid
ReadWritePaths=/var/lib/ongrid /var/log/ongrid

LimitNOFILE=65536
LimitNPROC=32768
MemoryMax=2G
CPUQuota=200%

Environment=GOMAXPROCS=auto
Environment=GOGC=100

[Install]
WantedBy=multi-user.target
