# deploy/install/systemd/prometheus.service
[Unit]
Description=Prometheus (ongrid metrics TSDB)
Documentation=https://prometheus.io/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ongrid-prometheus
Group=ongrid
ExecStart=/usr/local/bin/prometheus \
  --config.file=/etc/ongrid/prometheus/prometheus.yml \
  --storage.tsdb.path=/var/lib/ongrid-prometheus \
  --storage.tsdb.retention.time=15d \
  --web.listen-address=127.0.0.1:9090 \
  --web.enable-remote-write-receiver \
  --web.enable-lifecycle
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30s
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ongrid-prometheus

NoNewPrivileges=true
ProtectSystem=true
ProtectHome=true
PrivateTmp=true
StateDirectory=ongrid-prometheus
StateDirectoryMode=0755
ReadWritePaths=/var/lib/ongrid-prometheus
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
