# deploy/install/systemd/qdrant.service
[Unit]
Description=Qdrant (ongrid vector DB for RAG)
Documentation=https://qdrant.tech/documentation/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ongrid-qdrant
Group=ongrid
Environment=QDRANT__SERVICE__HTTP_PORT=6333
Environment=QDRANT__SERVICE__GRPC_PORT=6334
Environment=QDRANT__STORAGE__STORAGE_PATH=/var/lib/ongrid-qdrant/storage
Environment=QDRANT__STORAGE__SNAPSHOTS_PATH=/var/lib/ongrid-qdrant/snapshots
WorkingDirectory=/var/lib/ongrid-qdrant
Environment=QDRANT__LOG_LEVEL=INFO
ExecStart=/usr/local/bin/qdrant
KillMode=mixed
TimeoutStopSec=30s
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ongrid-qdrant

NoNewPrivileges=true
ProtectSystem=true
ProtectHome=true
PrivateTmp=true
StateDirectory=ongrid-qdrant
StateDirectoryMode=0755
ReadWritePaths=/var/lib/ongrid-qdrant
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
